# stale/mariadb/rootfs/etc/services.d/mysqld/run
#!/bin/sh -e
exec 2>&1
exec /usr/bin/mysqld --defaults-file=/data/conf/my.cnf --datadir=/data/db/mysql/ --basedir=/usr --plugin-dir=/usr/lib/mysql/plugin --user=mysql --console
